# src/main/resources/bin/profile.d/proxy.sh
#!/bin/sh -e

# To be expended by vagrant-up wrapper

