# src/main/resources/bin/profile.d/proxy.sh
#!/bin/sh -e

echo "Appliying proxy settings.."

# To be expanded by vagrant-up wrapper

